# config-templates/supervisord.conf
[unix_http_server]
file=/supervisor/supervisor.sock

[supervisord]
logfile=/supervisor/supervisord.log
pidfile=/supervisor/supervisord.pid
childlogdir=/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///supervisor/supervisor.sock

[program:laravel-worker]
process_name=%(program_name)s
command=php /project/site/artisan queue:work --tries=3
autostart=true
autorestart=true
numprocs=1